# crontab
0 0 * * * python /app/scripts/export.py
